# roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
User=node_exporter
ExecStart={{ node_exporter_bin_path }} {{ node_exporter_options }}

[Install]
WantedBy=multi-user.target
